# ansible/roles/stacks/docker/tasks/templates/docker-cleanup.timer.j2
[Unit]
Description=Regular cleanup of Docker images to free up disk space

[Timer]
OnCalendar=weekly
Unit=docker-cleanup.service

[Install]
WantedBy=timers.target
